# luci-app-bandix-plus/root/etc/uci-defaults/40_luci-bandix-plus
#!/bin/sh
# Register bandix-plus UCI changes with LuCI's apply mechanism so that
# saving the form automatically reloads /etc/init.d/bandix-plus.

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@bandix_plus[-1]
	add ucitrack bandix_plus
	set ucitrack.@bandix_plus[-1].init=bandix-plus
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache /tmp/luci-modulecache/* 2>/dev/null
exit 0
